# deploy/crypto-update.path
[Unit]
Description=Crypto Wallet Tracker — Deploy request watcher
After=multi-user.target

[Path]
PathExists=/var/lib/docker/volumes/crypto-wallet-tracker_wallet-data/_data/deploy/request.json

[Install]
WantedBy=multi-user.target
